# bot.conf
[program:testbot]
directory=/home/ubuntu/mukammal-bot-paid
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu/",USER="ubuntu"

stderr_logfile=/home/ubuntu/toonify_bot/logfile_err.log
stdout_logfile=/home/ubuntu/toonify_bot/logfile.log